# backlight-changer/backlight.service
[Unit]
Description=Change screen backlight

[Service]
ExecStart=/usr/sbin/backlight-changer.sh

[Install]
WantedBy=default.target
